# python-app/pkg/supervisor/python-app.ini
[program:python-app]
directory=/opt/python-app/src/
command=/opt/python-app/src/venv/bin/python3 run_python-app.py
user=root
autostart=true
autorestart=true
redirect_stderr=true

environment=CONFIG_FILE=/opt/python-app/config/python-app.yml
