# {{cookiecutter.project}}/pkg/supervisor/{{cookiecuttter.project_name}}.ini
[program:{{cookiecutter.project_name}}]
directory=/opt/{{cookiecutter.project_name}}/src/
command=/opt/{{cookiecutter.project_name}}/src/venv/bin/python3 run.py
user=root
autostart=true
autorestart=true
redirect_stderr=true

environment=CONFIG_FILE=/opt/{{cookiecutter.project_name}}/config/config.yaml
