# {{cookiecutter.project}}/pkg/supervisor/{{cookiecutter.project}}.ini
[program:{{cookiecutter.project}}]
directory=/opt/{{cookiecutter.project}}/src/
command=/opt/{{cookiecutter.project}}/src/venv/bin/python3 run.py
user=root
autostart=true
autorestart=true
redirect_stderr=true

environment=CONFIG_FILE=/opt/{{cookiecutter.project}}/config/config.yaml
